# docs/install/systemd/warden.target
[Unit]
Description=Warden
Wants=warden.scheduler.service warden.api.service